# installer/linux/scripts/fedora/vault_manager.service
[Unit]
Description=Vault Manager

[Service]
Type=simple
ExecStart=/usr/bin/maidsafe_vault_manager --start
 
[Install]
WantedBy=multi-user.target
